# luci-app-frpc/root/etc/init.d/frpc
#!/bin/sh /etc/rc.common
#
# Copyright 2019 Xingwang Liao <kuoruan@gmail.com>
# Licensed to the public under the MIT License.
#

START=99
USE_PROCD=1

NAME="frpc"
CONFIG_FOLDER="/var/etc/$NAME"

_log() {
	local level="$1" ; shift
	local msg="$@"
	logger -p "daemon.$level" -t "$NAME" "$msg"

	echo "[$level] $msg" >&2
}

_info() {
	_log "info" $@
}

_err() {
	_log "err" $@
}
# 配置生成已迁移到 Lua (Schema 驱动, toml_gen.lua)
# 由 /usr/lib/lua/frpc/gen.lua 处理 UCI → TOML 转换

frpc_scetion_validate() {
	uci_validate_section "$NAME" "frpc" "$1" \
		'enabled:bool:0' \
		'client_file:file:/usr/bin/frpc' \
		'server:uci("frpc", "@server")' \
		'run_user:string' \
		'enable_logging:bool:0' \
		'std_redirect:bool:0' \
		'log__to:string' \
		'log__level:or("trace", "debug", "info", "warn", "error")' \
		'log__maxDays:uinteger' \
		'log__disablePrintColor:or("true", "false")' \
		'loginFailExit:or("true", "false")' \
		'transport__poolCount:uinteger' \
		'user:string' \
		'transport__protocol:string' \
		'transport__quic__keepalivePeriod:integer' \
		'transport__quic__maxIdleTimeout:integer' \
		'transport__quic__maxIncomingStreams:uinteger' \
		'transport__proxyURL:string' \
		'transport__tls__enable:or("true", "false")' \
		'transport__tls__disableCustomTLSFirstByte:or("true", "false")' \
		'transport__tls__certFile:file' \
		'transport__tls__keyFile:file' \
		'transport__tls__trustedCaFile:file' \
		'transport__tls__serverName:string' \
		'dnsServer:host' \
		'natHoleStunServer:string' \
		'transport__heartbeatInterval:integer' \
		'transport__heartbeatTimeout:integer' \
		'com_extra_options:list(string)' \
		'webServer__addr:host' \
		'webServer__port:port' \
		'webServer__user:string' \
		'webServer__password:string'
}

server_section_validate() {
	uci_validate_section "$NAME" "server" "$1" \
		'alias:string' \
		'serverAddr:string' \
		'serverPort:port' \
		'auth__method:or("token", "oidc")' \
		'auth__token:string' \
		'transport__tcpMux:or("true", "false")' \
		'transport__tcpMuxKeepaliveInterval:uinteger'
}

client_file_validate() {
	local file="$1"

	test -f "$file" || return 1
	test -x "$file" || chmod 755 "$file"

	eval "$file" -h | grep -q "$NAME"
	return $?
}

create_config_file() {
	local config_file="$1"
	lua /usr/lib/lua/frpc/gen.lua "$NAME" "$config_file"
}
logfile_prepare() {
	if [ "x$enable_logging" != "x1" ];then
		rm -f /tmp/frpc_log_link.txt
		rm -f /tmp/frpc_std_redirect.log
		echo "未配置日志输出保存" > /tmp/frpc_log_link.txt
		return 1
	fi
	
	if [ "${std_redirect}" == "1" ];then
		if [ -z "${log__to}" -o "${log__to}" == "console" ];then
			log__to=/tmp/frpc_std_redirect.log
			true >"${log__to}"
		else
			mkdir -p "$(dirname "${log__to}")"
			touch "${log__to}"
			if [ "$?" != "0" ];then
				log__to=/tmp/frpc_std_redirect.log
				echo "【$(date +%Y-%m-%d_%H:%M:%S)】--->警告：配置的日志文件路径似乎无效/只读！" > "${log__to}"
			fi
		fi
	fi
	
	if [ -n "${log__to}" -a "${log__to}" != "console" ];then
		mkdir -p "$(dirname "${log__to}")"
		touch "${log__to}"
	fi
	
	if [ -n "$run_user" ]; then
		chmod 644 "$log__to"
		chown "$run_user" "$log__to"
	fi
	# 创建log文件读取链接，用于luci查看日志
	ln -sf "${log__to}" /tmp/frpc_log_link.txt
}
start_instance() {
	local section="$1"

	if ! frpc_scetion_validate "$section" ; then
		_err "Config validate failed."
		return 1
	fi

	if [ "x$enabled" != "x1" ] ; then
		_info "Instance \"$section\" disabled."
		return 1
	fi

	if [ -z "$client_file" ] || ( ! client_file_validate "$client_file" ) ; then
		_err "Client file not valid."
		return 1
	fi
	
	if [ -z "$server" ] || [ "$server" = "nil" ] ; then
		_err "No server selected for instance: \"$section\"."
		return 1
	elif ! server_section_validate "$server" ; then
		_err "Server config validate failed: \"$server\""
		return 1
	fi

	test -d "$CONFIG_FOLDER" || mkdir -p "$CONFIG_FOLDER"

	local config_file="$CONFIG_FOLDER/frpc.$section.toml"

	create_config_file "$config_file"

	if [ ! -f "$config_file" ] ; then
		_err "Could not create config file: \"$config_file\""
		return 1
	fi
	
	logfile_prepare

	procd_open_instance "$NAME.$section"
	if [ "${std_redirect}" == "1" ];then
		procd_set_param command "/bin/sh"
		procd_append_param command -c "exec $client_file -c $config_file >>\"${log__to}\" 2>&1"
	else
		procd_set_param command "$client_file"
		procd_append_param command -c "$config_file"
	fi
	procd_set_param respawn
	procd_set_param file "$config_file"

	if [ -n "$run_user" ] ; then
		procd_set_param user "$run_user"
	fi

	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger "$NAME"
}

start_service() {
	config_load "$NAME"
	config_foreach start_instance "frpc"
}
